# media-ingest/systemd/media-ingest.service
[Unit]
Description=Media Ingest (app + worker + Caddy)
After=network.target tailscaled.service
Wants=tailscaled.service

# Edit User, WorkingDirectory, and ExecStart to match your deployment
[Service]
Type=simple
User=<your-username>
WorkingDirectory=/opt/media-ingest
ExecStart=/opt/media-ingest/run.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
